# runit/backplane/run
#!/bin/sh

if [ -n "$BACKPLANE_TOKEN" ]
then
		/usr/bin/backplane connect
else
		sleep 999d
fi
